# tickstack/influxd.ini
[program:influxd]
command = /usr/bin/influxd run
