# alerta/docker/supervisord.d/uwsgi.ini
[program:uwsgi]
command = /venv/bin/uwsgi --ini /etc/uwsgi.ini
program_name = uwsgi
