# elk_stack/logstash/logstash.ini
[program:logstash]
command = /usr/share/logstash/bin/logstash  --path.settings /etc/logstash
autostart = true
user = logstash
